# crontab
* * * * * root /usr/bin/make -f /code/Makefile all >> /var/log/download-data.log 2>&1
# Empty line.
