# dashboard/crontab
0 1 * * * root cd /code && /usr/bin/make all >> /var/log/download-data.log 2>&1
# Empty line.
